# build/systemd/aurora-bot.dbus.service
[Unit]
Description=Aurora Bot | Dbus API service
Requires=default.target
After=default.target

[Service]
Environment=DISPLAY=:1
Type=exec
KillMode=process
TimeoutSec=5min
Restart=no
ExecStart=/usr/local/bin/aurora-bot svc --dbus
WorkingDirectory=%h

[Install]
WantedBy=default.target
