# build/systemd/aurora-bot.client.service
[Unit]
Description=Aurora Bot | Telegram client service
Requires=default.target
After=default.target

[Service]
Environment=DISPLAY=:1
Type=exec
KillMode=process
TimeoutSec=5min
Restart=no
ExecStart=/usr/bin/aurora-bot svc --connect
WorkingDirectory=%h

[Install]
WantedBy=default.target
